# dockertunnel.service
[Unit]
Description=Maintain a tunnel with the docker host of the vps

[Service]
ExecStart=/usr/local/bin/dockertunnel.sh
